# setup/docker/crontab
*/5 * * * * root /usr/bin/python /opt/nest-datagraph/backend/poller.py
